# momcontacts.service
[Unit]
Description=Mom Contacts Management System
After=network.target

[Service]
Type=simple
User=samadjei
WorkingDirectory=/Users/samadjei/MomContacts
ExecStart=/usr/bin/python3 /Users/samadjei/MomContacts/app.py
Restart=always
RestartSec=10
Environment=FLASK_ENV=production

[Install]
WantedBy=multi-user.target 